# system_files_overrides/gdx/usr/lib/systemd/system/bluefin-lts-migration.timer
[Unit]
Description=Daily trigger for Bluefin LTS image migration
Documentation=https://projectbluefin.io

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
Unit=bluefin-lts-migration.service

[Install]
WantedBy=timers.target
