# system_files/usr/lib/systemd/system/ublue-countme.timer
[Unit]
Description=Weekly update of dnf countme

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
