# system_files/usr/lib/systemd/system/ublue-countme.service
[Unit]
Description=Enable countme for dnf

[Service]
Type=oneshot
ExecStart=/usr/bin/ublue-countme
